# pi-setup/systemd/hanryxvault-barcode.service
[Unit]
Description=HanryxVault QR Scan Hub (USB + Bluetooth HID, broadcasts to all apps)
After=hanryxvault.service bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hanryxvault
Environment=HANRYX_DIR=/opt/hanryxvault
Environment=SCAN_HUB_PORT=8765
ExecStart=/opt/hanryxvault/venv/bin/python3 /opt/hanryxvault/barcode_daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hanryxvault-scan-hub

[Install]
WantedBy=multi-user.target
