# pi-setup/systemd/hanryx-heal.timer
[Unit]
Description=HanryxVault — run recover.sh every 2 minutes
Documentation=file:///home/ngansen/Hanryx-Vault-POS/pi-setup/recover.sh

[Timer]
# Run 30s after boot, then every 60 seconds.
# 60s cadence catches WiFi flaps, dead Chromium, container crashes
# within a single trade-show transaction window. recover.sh is a no-op
# when everything is healthy so this is cheap.
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Persistent=false
Unit=hanryx-heal.service

[Install]
WantedBy=timers.target
