# pi-setup/systemd/hanryxvault-satellite-sync.service
[Unit]
Description=HanryxVault Satellite Monitor (offline sales storage + auto-sync on reconnect)
After=network-online.target wg-quick@wg0.service hanryxvault.service
Wants=network-online.target wg-quick@wg0.service
Requires=hanryxvault.service

[Service]
Type=simple
User=hanryxvault
WorkingDirectory=/opt/hanryxvault
# Brief delay to let WireGuard tunnel establish before first connectivity check
ExecStartPre=/bin/sleep 8
ExecStart=/opt/hanryxvault/venv/bin/python3 /opt/hanryxvault/satellite_monitor.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hanryxvault-sync
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
