# pi-setup/kiosk/hanryxvault-kiosk-restart.service
[Unit]
Description=HanryxVault — Nightly kiosk restart (clears chromium memory leaks)
Wants=hanryxvault-kiosk.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart hanryxvault-kiosk.service
